# rpicontrol.service
[Unit]
Description=RPi Control GUI - Load Cell and Motor Control System
After=graphical-session.target network.target
Wants=graphical-session.target

[Service]
Type=simple
User=%i
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%i
Environment=PYTHONPATH=%h/projects/Rpi_force_control
Environment=PYTHONUNBUFFERED=1
Environment=QT_X11_NO_MITSHM=1
# Wait for display to be ready and activate virtual environment
ExecStart=/bin/bash -c 'cd %h/projects/Rpi_force_control && python rpi_control.py'
WorkingDirectory=%h/projects/Rpi_force_control
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=120
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStartSec=30
TimeoutStopSec=15

[Install]
WantedBy=default.target graphical-session.target
